# systemd/etf-data-update.timer
[Unit]
Description=Run ETF DB aggregation update at 20:00 China time (12:00 UTC)

[Timer]
OnCalendar=*-*-* 12:00:00 UTC
Persistent=true
Unit=etf-data-update.service

[Install]
WantedBy=timers.target
