# systemd/etf-data-update.timer
[Unit]
Description=Run ETF DB aggregation update at 18:00 China time (10:00 UTC)

[Timer]
OnCalendar=*-*-* 10:00:00 UTC
Persistent=true
Unit=etf-data-update.service

[Install]
WantedBy=timers.target
